# etc/profile.d/gvm.sh
# shellcheck disable=SC2148
export PATH="$PATH:/opt/gvm/bin:/opt/gvm/sbin:/opt/gvm/.local/bin"
