# etc/systemd/system/gvmd-feedsync.service
[Unit]
Description=Synchronize GVM feeds
After=network.target prepare-gvm.service

[Service]
Type=oneshot
ExecStart=/opt/gvm/sbin/feedsync.sh
