# etc/systemd/system/gvmd-feedsync.timer
[Unit]
Description=Synchronize GVM feeds

[Timer]
OnCalendar=*-*-* 12:00:00
RandomizedDelaySec=7200

[Install]
WantedBy=multi-user.target
