# etc/systemd/system/update-gvm-targets.timer
[Unit]
Description=Synchronize GVM feeds

[Timer]
OnCalendar=*-*-* 17:00:00

[Install]
WantedBy=multi-user.target
